# source/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:api]
command=python main.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api_err.log
environment=
    UPLOAD_DIR="/tmp/cubers_uploads",
    MAX_UPLOAD_MB="200",
    INTENSITY_THRESHOLD="2.0",
    MIN_PAUSE_DURATION="0.5",
    PORT="8000"

[program:streamlit]
command=streamlit run app.py
    --server.port=8501
    --server.address=0.0.0.0
    --server.headless=true
    --server.maxUploadSize=200
    --browser.gatherUsageStats=false
    --theme.base=dark
    --theme.primaryColor=#818cf8
    --theme.backgroundColor=#0d0d1a
    --theme.secondaryBackgroundColor=#111827
    --theme.textColor=#e5e7eb
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/streamlit.log
stderr_logfile=/var/log/supervisor/streamlit_err.log
environment=
    API_URL="http://localhost:8000",
    MAX_UPLOAD_MB="200"
